# roles/clamav/files/clamd
#!/bin/ksh

daemon="/usr/local/sbin/clamd"
daemon_timeout="120"

. /etc/rc.d/rc.subr

rc_reload=NO
rc_bg=YES

rc_cmd $1
